# files/etc/systemd/system/universal-lite-first-boot.service
[Unit]
Description=First-boot configuration for installer-deployed Universal-Lite
After=local-fs.target
Before=greetd.service swap.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-first-boot
TimeoutStartSec=600
ConditionPathExists=/var/lib/universal-lite/install-config.json
ConditionPathExists=!/var/lib/universal-lite/setup-done

[Install]
WantedBy=multi-user.target
