# files/etc/systemd/system/universal-lite-flatpak-setup.service
[Unit]
Description=Install Flatpak apps selected during setup wizard
Wants=NetworkManager-wait-online.service network-online.target
After=NetworkManager-wait-online.service network-online.target
ConditionPathExists=!/var/lib/universal-lite/flatpak-setup.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-flatpak-setup
RemainAfterExit=yes
TimeoutStartSec=1800

[Install]
WantedBy=graphical.target
