# files/etc/systemd/system/universal-lite-flatpak-setup.service
[Unit]
Description=Install and keep Flatpak apps up to date
# Deliberately NOT Before=greetd.service: login should come up
# immediately and flatpak installs / updates run in the background.
# Deliberately no NetworkManager-wait-online dep: the script polls
# dl.flathub.org DNS itself, which is both more accurate and doesn't
# delay login on offline boots (wait-online times out at 30s).
After=universal-lite-first-boot.service NetworkManager.service
# Only run after the wizard has completed (setup-done exists). The
# script itself decides whether to install (stamp absent) or update
# (stamp present) based on /var/lib/universal-lite/flatpak-setup.done,
# so the service runs on every boot — which is how the user gets
# keep-current Flatpaks for free.
ConditionPathExists=/var/lib/universal-lite/setup-done

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-flatpak-setup
RemainAfterExit=yes
TimeoutStartSec=1800
StandardOutput=journal+console

[Install]
WantedBy=graphical.target
