# files/etc/systemd/system/universal-lite-flatpak-install.service
[Unit]
Description=First-boot Flatpak install
# Runs BEFORE greetd so the user logs into a complete desktop with
# their selected apps already present, not an empty one with
# background downloads in flight. Hard-timeout safety valve below
# prevents a flaky network from locking the user out forever.
Before=greetd.service
After=universal-lite-first-boot.service NetworkManager.service
# Wizard must have completed AND we must not already be past the
# install phase. Once the stamp exists, the -update.service unit
# takes over and runs in the background.
ConditionPathExists=/var/lib/universal-lite/setup-done
ConditionPathExists=!/var/lib/universal-lite/flatpak-setup.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-flatpak-setup
RemainAfterExit=yes
# 10 minutes is more than enough for the default four apps on even
# VM-speed networks. If it's still running at that point something
# is wrong (permanent Flathub outage, etc.) — fall through to login
# rather than trap the user. The script writes no stamp on failure
# and next boot resumes where it left off (already-installed apps
# are skipped via `flatpak list`).
TimeoutStartSec=600
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=graphical.target
