# files/etc/systemd/system/universal-lite-flatpak-setup.service
[Unit]
Description=Universal-Lite Flatpak setup (Flathub + default apps)
After=network-online.target
Wants=network-online.target
ConditionPathIsReadWrite=/var/lib/flatpak

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-flatpak-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
