# files/etc/systemd/system/universal-lite-flatpak-install.service
[Unit]
Description=First-boot Flatpak install
# Runs in parallel with greetd. The greeter watches the progress
# file (/run/universal-lite/flatpak-progress) and shows a "Finishing
# setup" overlay in place of the login form until the stamp file
# appears, so login is blocked from the user's perspective but the
# greeter itself can come up instantly with actual progress
# feedback instead of a silent Plymouth splash.
After=universal-lite-first-boot.service NetworkManager.service
# Wizard must have completed AND we must not already be past the
# install phase. Once the stamp exists, the -update.service unit
# takes over and runs in the background.
ConditionPathExists=/var/lib/universal-lite/setup-done
ConditionPathExists=!/var/lib/universal-lite/flatpak-setup.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-flatpak-setup
RemainAfterExit=yes
# 10 minutes is more than enough for the default four apps on even
# VM-speed networks. If it's still running at that point something
# is wrong (permanent Flathub outage, etc.) — fall through to login
# rather than trap the user. The script writes no stamp on failure
# and next boot resumes where it left off (already-installed apps
# are skipped via `flatpak list`).
TimeoutStartSec=600
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=graphical.target
