# files/usr/lib/systemd/system/user@.service.d/10-universal-lite-oomd.conf
[Service]
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=55%
