# files/etc/systemd/system/greetd.service.d/no-plymouth-wait.conf
[Unit]
# Remove plymouth-quit-wait dependency so greetd starts regardless
# of plymouth state. Prevents boot hang when cage/greeter fails to
# signal plymouth.
After=
After=systemd-user-sessions.service
