# files/etc/systemd/system/uupd.service.d/10-first-boot-guard.conf
[Unit]
# Don't run uupd until the wizard + first-boot setup is complete.
# The Persistent=true timer fires immediately on first boot, but
# network isn't ready and first-boot services need resources.
ConditionPathExists=/var/lib/universal-lite/setup-done
After=network-online.target
Wants=network-online.target
