# files/usr/lib/systemd/user/universal-lite-chrome-oom-cleanup.service
[Unit]
Description=Clean up remaining Chrome Flatpak scopes after OOM

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-chrome-oom-cleanup
