# files/usr/lib/systemd/user/universal-lite-nightlight.service
[Unit]
Description=Reconcile night-light state with current time

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-nightlight
