# files/usr/lib/systemd/system/universal-lite-chrome-early-oom.service
[Unit]
Description=Universal-Lite Chrome early-OOM safety valve
Documentation=file:/usr/share/doc/universal-lite/chrome-early-oom.md
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/libexec/universal-lite-chrome-early-oom
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
