# files/usr/lib/systemd/system/universal-lite-nm-start.service
[Unit]
# Imperative NetworkManager starter. Bypasses whatever is suppressing
# the declarative Wants= on multi-user.target and greetd.service.d/
# 10-ensure-networkmanager.conf — on failing boots we can see
# `systemctl show greetd -p Wants` return Wants=NetworkManager.service
# and `systemctl cat NetworkManager` loaded+enabled, yet systemd never
# queues NM for activation.
#
# A/B test across 6 reboots on commit 9266b19 (this service removed,
# units relocated to /usr/lib/systemd/system/) still produced 2/6
# failures — so the relocation wasn't the full fix. Re-landing this
# as the pragmatic mitigation while upstream root cause is still open.
#
# This service uses ExecStart= to call systemctl directly, which issues
# a fresh start job for NetworkManager at runtime. That job doesn't
# go through the same transaction resolution that was silently dropping
# the Wants=.

Description=Force-start NetworkManager (workaround for dropped Wants=)
After=dbus.service
Before=greetd.service
Wants=dbus.service

[Service]
Type=oneshot
# --no-block so we don't wait for NM to fully settle. The actual NM
# start is backgrounded; this service is considered done as soon as
# the job is queued.
ExecStart=/usr/bin/systemctl --no-block start NetworkManager.service
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
