# files/usr/lib/systemd/system/universal-lite-zswap.service
[Unit]
Description=Configure zswap for disk swap
DefaultDependencies=no
After=universal-lite-swap-init.service
Before=swap.target
Before=universal-lite-encrypted-swap.service
ConditionPathExists=/var/swap

[Service]
Type=oneshot
ExecStart=-/bin/sh -c 'echo 1 > /sys/module/zswap/parameters/enabled'
ExecStart=-/bin/sh -c 'echo 25 > /sys/module/zswap/parameters/max_pool_percent'
ExecStart=-/bin/sh -c 'echo zstd > /sys/module/zswap/parameters/compressor'
RemainAfterExit=yes

# --- Hardening ---
# Writes /sys/module/zswap/parameters/* — explicitly leaves ProtectKernelTunables
# OFF (default) and whitelists the sysfs path via ReadWritePaths.
ProtectSystem=strict
ProtectHome=yes
# PrivateTmp omitted: runs before /tmp is mountable (swap.target).
ReadWritePaths=/sys/module/zswap/parameters
NoNewPrivileges=yes
RestrictSUIDSGID=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelLogs=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
PrivateNetwork=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@clock @debug @module @reboot @obsolete @cpu-emulation @swap @mount
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=CAP_DAC_OVERRIDE
AmbientCapabilities=

[Install]
WantedBy=swap.target
