# files/usr/lib/systemd/system/universal-lite-flatpak-install.service
[Unit]
Description=First-boot Flatpak install
# Runs in parallel with greetd. The greeter watches the progress
# file (/run/universal-lite/flatpak-progress) and shows a "Finishing
# setup" overlay in place of the login form until the stamp file
# appears, so login is blocked from the user's perspective but the
# greeter itself can come up instantly with actual progress
# feedback instead of a silent Plymouth splash.
After=universal-lite-first-boot.service NetworkManager.service
# Wizard must have completed AND we must not already be past the
# install phase. Once the stamp exists, the -update.service unit
# takes over and runs in the background.
ConditionPathExists=/var/lib/universal-lite/setup-done
ConditionPathExists=!/var/lib/universal-lite/flatpak-setup.done

[Service]
# Type=exec so graphical.target reaches active as soon as bash is
# exec'd, not when every app is installed. The greeter overlay is
# driven by the stamp file + progress file (not by systemd state),
# so the UX is unchanged - the greeter dismisses the "Finishing
# setup" page the moment the script writes the stamp, regardless of
# whether systemd is still tracking the unit as running.
Type=exec
ExecStart=/usr/libexec/universal-lite-flatpak-setup
StandardOutput=journal+console
StandardError=journal+console

# --- Hardening ---
# flatpak install is a heavyweight operation: downloads from dl.flathub.org,
# writes to /var/lib/flatpak, /var/tmp, creates OSTree refs, may spawn bwrap
# for triggers. Skip NoNewPrivileges (flatpak's sandbox helpers can legitimately
# need setuid bwrap on non-userns kernels) and ProtectKernelTunables (some
# flatpak runtimes poke sysfs via triggers).
ProtectSystem=strict
# flatpak writes /var/lib/flatpak, /var/tmp, and we write /run/universal-lite/
ReadWritePaths=/var /run
# Not ProtectHome=yes: flatpak's system installs don't touch /home, but
# user-mode triggers may. Use read-only as a safe middle ground.
ProtectHome=read-only
PrivateTmp=yes
RestrictSUIDSGID=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native
# flatpak uses TCP (HTTPS to flathub), AF_UNIX (dbus, journal), AF_NETLINK.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=graphical.target
