# files/usr/lib/systemd/system/universal-lite-first-boot.service
[Unit]
Description=First-boot configuration for installer-deployed Universal-Lite
# Before=swap.target was the cycle source: combined with the default
# After=basic.target it forced sysinit→basic→first-boot→swap ordering,
# which collided with swap.target's placement in the boot sequence.
# First-boot writes swap-size and calls the swap scripts directly
# (see /usr/libexec/universal-lite-first-boot), so blocking swap.target
# from systemd's side is unnecessary.
After=local-fs.target
Before=greetd.service
ConditionPathExists=/var/lib/universal-lite/install-config.json
ConditionPathExists=!/var/lib/universal-lite/setup-done

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-first-boot
TimeoutStartSec=600

# --- Hardening ---
# Intentionally light: this script calls systemctl enable, modprobe zswap,
# sysctl -w, writes to /sys/module/zswap/parameters, touches /var/lib/universal-lite,
# creates /var/swap, and execs other libexec scripts. Too much to lock down
# without breaking it. Apply only the directives that are safely compatible:
#   - ProtectSystem=strict + ReadWritePaths covers /var writes.
#   - NoNewPrivileges=yes is fine (no setuid execs involved).
#   - Skip ProtectKernelModules (modprobe), ProtectKernelTunables (sysctl/sysfs).
ProtectSystem=strict
ReadWritePaths=/var /run /sys/module/zswap/parameters /proc/sys/vm
NoNewPrivileges=yes
RestrictSUIDSGID=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
