# files/usr/lib/systemd/system/universal-lite-swap-init.service
[Unit]
Description=Create encrypted swap file
# After=local-fs.target caused an ordering cycle: Fedora's tmp.mount
# has After=swap.target, and local-fs.target pulls in tmp.mount, so
# local-fs.target ended up transitively after swap.target -- which
# can't also be before it. RequiresMountsFor=/var pulls in just the
# specific mount unit the swap file needs, no target transits.
DefaultDependencies=no
RequiresMountsFor=/var
Before=swap.target universal-lite-encrypted-swap.service
ConditionPathExists=!/var/swap
ConditionPathExists=/var/lib/universal-lite/swap-size

[Service]
Type=oneshot
ExecStart=/usr/libexec/universal-lite-swap-init
RemainAfterExit=yes

# --- Hardening ---
# Writes /var/swap (only on first boot; guarded by ConditionPathExists=!/var/swap).
# Uses chattr +C (CAP_LINUX_IMMUTABLE) on btrfs, fallocate or dd, chmod.
# Reads /dev/zero and /var/lib/universal-lite/swap-size.
ProtectSystem=strict
ProtectHome=yes
# PrivateTmp omitted: runs before /tmp is mountable (same cycle as above).
ReadWritePaths=/var
NoNewPrivileges=yes
RestrictSUIDSGID=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelLogs=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
PrivateNetwork=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@clock @debug @module @reboot @obsolete @cpu-emulation @swap
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=CAP_LINUX_IMMUTABLE CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE
AmbientCapabilities=

[Install]
WantedBy=swap.target
