# files/usr/lib/systemd/system/universal-lite-flatpak-install.service
[Unit]
Description=First-boot Flatpak install
# Runs in parallel with greetd. The greeter watches the progress
# file (/run/universal-lite/flatpak-progress) and shows a "Finishing
# setup" overlay in place of the login form until the stamp file
# appears, so login is blocked from the user's perspective but the
# greeter itself can come up instantly with actual progress
# feedback instead of a silent Plymouth splash.
After=universal-lite-first-boot.service NetworkManager.service
# Wizard must have completed AND we must not already be past the
# install phase. Once the stamp exists, the -update.service unit
# takes over and runs in the background.
ConditionPathExists=/var/lib/universal-lite/setup-done
ConditionPathExists=!/var/lib/universal-lite/flatpak-setup.done

[Service]
# Type=exec so graphical.target reaches active as soon as bash is
# exec'd, not when every app is installed. The greeter overlay is
# driven by the stamp file + progress file (not by systemd state),
# so the UX is unchanged - the greeter dismisses the "Finishing
# setup" page the moment the script writes the stamp, regardless of
# whether systemd is still tracking the unit as running.
Type=exec
ExecStart=/usr/libexec/universal-lite-flatpak-setup
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=graphical.target
