# files/usr/lib/systemd/system/universal-lite-encrypted-swap.service
[Unit]
Description=Encrypted swap
# DefaultDependencies=no matches the other swap services (swap-init,
# zswap). Without it, the implicit After=sysinit.target combines with
# Before=swap.target + WantedBy=swap.target to force sysinit-before-swap,
# which creates the ordering cycle systemd was logging on every boot.
DefaultDependencies=no
RequiresMountsFor=/var/swap
After=universal-lite-swap-init.service
Before=swap.target
ConditionPathExists=/var/swap

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/universal-lite-encrypted-swap start
ExecStop=/usr/libexec/universal-lite-encrypted-swap stop

[Install]
WantedBy=swap.target
