# files/usr/lib/systemd/system/universal-lite-encrypted-swap.service
[Unit]
Description=Encrypted swap
# DefaultDependencies=no matches the other swap services (swap-init,
# zswap). Without it, the implicit After=sysinit.target combines with
# Before=swap.target + WantedBy=swap.target to force sysinit-before-swap,
# which creates the ordering cycle systemd was logging on every boot.
DefaultDependencies=no
RequiresMountsFor=/var/swap
After=universal-lite-swap-init.service
Before=swap.target
ConditionPathExists=/var/swap

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/universal-lite-encrypted-swap start
ExecStop=/usr/libexec/universal-lite-encrypted-swap stop

# --- Hardening ---
# Needs block device + dm-crypt + swapon + /dev/mapper/control + /dev/loopN +
# /dev/urandom. Keeps CAP_SYS_ADMIN (loop/dm-crypt/swapon) — everything else
# can go. Writes runtime loop state to /run and opens /var/swap read/write
# as the loop backing file.
ProtectSystem=strict
ProtectHome=yes
# PrivateTmp omitted: swap.target runs before /tmp is mountable, so
# pulling in systemd-tmpfiles-setup would create an ordering cycle.
ReadWritePaths=/run /var/swap
NoNewPrivileges=yes
RestrictSUIDSGID=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelLogs=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native
# @swap + @mount needed for swapon + dm-crypt ioctls. @raw-io not needed —
# cryptsetup uses the device-mapper ioctl interface.
SystemCallFilter=@system-service @swap @mount
SystemCallFilter=~@clock @debug @module @reboot @obsolete @cpu-emulation
# dm-crypt/cryptsetup open AF_NETLINK for uevents; AF_UNIX for journald.
RestrictAddressFamilies=AF_UNIX AF_NETLINK
CapabilityBoundingSet=CAP_SYS_ADMIN
AmbientCapabilities=

[Install]
WantedBy=swap.target
