# files/usr/lib/systemd/system/greetd.service.d/10-ensure-networkmanager.conf
# NetworkManager.service is WantedBy=multi-user.target and our build.sh
# explicitly `systemctl enable`s it, but field reports show that after
# a fresh bootc-installed boot it sometimes stays inactive (dead) with
# zero journal entries — systemd reached multi-user.target but never
# activated the Want. Unclear whether it's a first-boot generator-cache
# race or an interaction with NM's upstream unit conditions; either
# way, the user ends up at a login screen with no network.
#
# greetd fires reliably on every graphical boot, so pulling NM in via
# a Wants= + After= on greetd re-asserts the start request on every
# boot. Wants= (not Requires=) so greetd still comes up even if NM
# genuinely cannot start, rather than leaving the user stranded with
# no way to log in and debug.

[Unit]
Wants=NetworkManager.service
After=NetworkManager.service
