# files/usr/lib/systemd/system/greetd.service.d/10-greeter-setup.conf
[Service]
ExecStartPre=/usr/libexec/universal-lite-greeter-setup
